# deploy/wa-server.service
[Unit]
Description=Copa no Zap - servidor WhatsApp (Baileys)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=copa
WorkingDirectory=/opt/copa-do-mundo/wa-server
EnvironmentFile=/opt/copa-do-mundo/.env
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
